# build/papermc-server.conf
[program:papermc-server]
user = nobody
command = /home/nobody/start.sh
umask = 000
